# packaging/systemd/claude-contextd.service.in
[Unit]
Description=Claude Context daemon
Documentation=https://github.com/agoodkind/claude-context-go
After=network.target

[Service]
ExecStart=@@BIN_PATH@@
Restart=on-failure
RestartSec=2
Environment=HOME=@@HOME@@

[Install]
WantedBy=default.target
